# deb/debian/postinst
#!/bin/sh
# postinst script for emqttd
#
# see: dh_installdeb(1)

set -e

# install startup script
update-rc.d emqttd defaults >/dev/null

# create group
if ! getent group emqtt >/dev/null; then
        addgroup --system emqtt
fi

# create user
if ! getent passwd emqtt >/dev/null; then
        adduser --ingroup emqtt \
                --home /var/lib/emqttd \
                --disabled-password \
                --system --shell /bin/bash --no-create-home \
                --gecos "emqttd user" emqtt
fi

for i in lib run log; do
    chown -R emqtt:emqtt /var/$i/emqttd
done

chown -R emqtt:emqtt /usr/lib/emqttd
chown -R emqtt:emqtt /etc/emqttd

chmod 0755 /var/run/emqttd /etc/emqttd
chmod 0644 /etc/emqttd/*
chmod -R +X /etc/emqttd
chmod 0755 /usr/lib/emqttd/bin/nodetool

case "$1" in
    configure)
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
